# packer/files/puma.service
[Unit]
Description=Puma server
[Service]
WorkingDirectory=/home/morgul/reddit
ExecStart=/usr/local/bin/puma
[Install]
WantedBy=multi-user.target